# .config/systemd/user/battery-check.service
[Service]
#Type=simple
ExecStart=%h/.local/bin/scripts/battery-check.sh
