# .config/systemd/user/battery-check.timer
[Unit]
Description=Notify if the battery is running low

[Timer]
OnCalendar=*-*-* *:*:00
Unit=battery-check.service
